# ansible/roles/supervisor/templates/supervisor.conf.j2
[program:blog]
command=/home/deploy/app/start.sh
user=deploy
autostart=true
autorestart=true
stderr_logfile=/var/log/blog.err.log
stdout_logfile=/var/log/blog.out.log
